# ansible/roles/locust/templates/locust-master.service.j2
[Unit]
Description=Locust load testing service

[Service]
ExecStart=/opt/locust/venv/bin/locust --locustfile=/opt/locust/Locustfiles/locustfile.py --master --host=https://www.mir.de
WorkingDirectory=/opt/locust
User=nobody
Restart=always

[Install]
WantedBy=multi-user.target

